# homarr.container
[Container]
Image=ghcr.io/homarr-labs/homarr:v1.22.0
Volume=homarr-configs:/app/data/configs:Z
Volume=homarr-icons:/app/public/icons:Z
Volume=homarr-data:/app/data:Z
PublishPort=8082:7575

[Service]
Restart=on-failure

[Install]
WantedBy=default.target
